# debinstall/deb-src/sysroot/etc/init.d/baby-monitor
#!/bin/sh

### BEGIN INIT INFO
# Provides:        baby-monitor
# Required-Start:  $network $remote_fs $syslog
# Required-Stop:   $network $remote_fs $syslog
# Default-Start:   2 3 4 5
# Default-Stop:    0 1 6
# Short-Description: Start baby-monitor daemon
### END INIT INFO


PATH="/opt/baby-monitor/:/usr/local/bin:/usr/bin:${PATH}"

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting baby-monitor"
    cd "/opt/baby-monitor"
    su baby-monitor -c "npm start"
    ;;
  stop)
    echo "Stopping baby-monitor"
    cd "/opt/baby-monitor"
    su baby-monitor -c "npm stop"
    ;;
  *)
    echo "Usage: /etc/init.d/baby-monitor {start|stop}"
    exit 1
    ;;
esac

exit 0
